# infra/roles/app/files/app.service
[Unit]
Description=app

[Service]
Type=simple

ExecStart=/usr/bin/python /home/vagrant/app.py

[Install]
WantedBy=multi-user.target
